# trunk/larch/profiles/xmini/rootoverlay/etc/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#

# Set up automatically logged in user
if [ -f /.livesys/autologin ]; then
    cp /.livesys/autologin /tmp/newuser
fi

